# packaging/debian/DEBIAN/postinst
#!/bin/sh
set -e
# Source debconf library.
. /usr/share/debconf/confmodule
update-rc.d visibility defaults 90 >/dev/null
# start visibility
invoke-rc.d visibility start
db_stop
