# packaging/debian/DEBIAN/preinst
#!/bin/bash
set -e
# stop django server
if [ -f /etc/init.d/visibility ]
then
  invoke-rc.d django stop
fi
