# packaging/debian/DEBIAN/prerm
#!/bin/bash
set -e
 
# stop visibility (django) server
invoke-rc.d visibility stop
